# files/wpn_daily.cron
#!/usr/bin/env bash
printf "Searching for WordPress Networks.\n"
find -L /srv -maxdepth 1 -mindepth 1 -type d 2>/dev/null | while read -r SITEPATH ; do
  # Check for settings.php
  stat ${SITEPATH}/wp/wp-config.php &>/dev/null
  if [ "$?" -eq "0" ]; then
    printf "Backing up ${SITEPATH}\n"
    /opt/wpn/bin/wpn_snapshot.sh ${SITEPATH}
  fi
done

printf "Done.\n
